# dhcp/image/root/etc/services.d/dhcpd/run
#!/usr/bin/with-contenv bash

/usr/sbin/dhcpd -6 -f -d --no-pid -cf /etc/dhcp/dhcpd6.conf
